# stream-bot.service
[Unit]
Description=Khdiamond Telegram Stream Bot
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=ubuntu
# Replace /home/ubuntu/stream-bot-with-go with your actual path
WorkingDirectory=/home/ubuntu/stream-bot-with-go
EnvironmentFile=/home/ubuntu/stream-bot-with-go/.env
ExecStart=/home/ubuntu/stream-bot-with-go/stream-bot
Restart=always
RestartSec=5
# Ensure logs are visible in journalctl
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
